# hassio_addon/rootfs/etc/services.d/zigbee-sniff/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting zigbee sniffing"

while true; do
	interval="$(bashio::config capture_interval 60)"
	timeout $interval /whsniff-1.3/whsniff -c 14 > /tmp/.capture.pcap
	ret=$?
	if [ "$ret" != 0 ]; then
		bashio::log.warning "Stopping zigbee sniffing"
		exit $ret
	fi
	if [ "$(bashio::config save_capture 'false')" == 'true' ]; then
		dest="/share/solarcity_sniffer/capture-$(date +%Y%m%d_%H%M%S).pcap"
		bashio::log.info "Saving a segment of length $(wc -c /tmp/.capture.pcap | awk '{print $1}') to $dest"
		[ -d /share ] && [ ! -d /share/solarcity_sniffer ] && mkdir /share/solarcity_sniffer
		[ -d /share/solarcity_sniffer ] && cp /tmp/.capture.pcap "$dest" || true
	fi
	mv /tmp/.capture.pcap /tmp/capture.pcap
done

